# packages/clairview-boraq/clairview_boraq-5.0.4-py3-none-any.whl/boraq/config/templates/systemd/clairview-boraq-clairview-default-worker.service
[Unit]
Description="{{ boraq_name }}-clairview-default-worker %I"
PartOf={{ boraq_name }}-workers.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ boraq_cmd }} worker --queue default
StandardOutput=file:{{ boraq_dir }}/logs/worker.log
StandardError=file:{{ boraq_dir }}/logs/worker.error.log
WorkingDirectory={{ boraq_dir }}
