# packages/clairview-boraq/clairview_boraq-5.0.0.dev0-py3-none-any.whl/boraq/config/templates/systemd/clairview-boraq-redis.target
[Unit]
After=network.target
Wants={{ boraq_name }}-redis-cache.service {{ boraq_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
