# packages/clairview-boraq/clairview_boraq-5.0.0.dev0.tar.gz/clairview_boraq-5.0.0.dev0/boraq/config/templates/systemd/clairview-boraq-node-socketio.service
[Unit]
After={{ boraq_name }}-clairview-web.service
Description="{{ boraq_name }}-node-socketio"
PartOf={{ boraq_name }}-web.target

[Service]
User={{ user }}
Group={{ user }}
Restart=always
ExecStart={{ node }} {{ boraq_dir }}/apps/clairview/socketio.js
StandardOutput=file:{{ boraq_dir }}/logs/node-socketio.log
StandardError=file:{{ boraq_dir }}/logs/node-socketio.error.log
WorkingDirectory={{ boraq_dir }}
