# packages/clairview-boraq/clairview_boraq-5.0.0.dev0.tar.gz/clairview_boraq-5.0.0.dev0/boraq/config/templates/systemd/clairview-boraq-web.target
[Unit]
After=network.target
Wants={{ boraq_name }}-clairview-web.service {{ boraq_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
