# packages/clairview-boraq/clairview_boraq-5.0.4.tar.gz/clairview_boraq-5.0.4/boraq/config/templates/systemd/clairview-boraq.target
[Unit]
After=network.target
Requires={{ boraq_name }}-web.target {{ boraq_name }}-workers.target {{ boraq_name }}-redis.target

[Install]
WantedBy=multi-user.target
