# packages/aprilcmd/aprilcmd-0.1.10-py3-none-any.whl/aprilcmd/shelltool/etc/cron.hourly/april.h.sh
#!/bin/bash
echo "ycmd"
/usr/bin/ycmd.py host._self.client up ignore=window
